# runtime/daemon
#!/sbin/openrc-run

command_background="true"
command="/usr/bin/env"
command_args="NODE_ENV=production DYNO_PORT=80 DYNO_DIR=/mnt/system/calc node --watch /mnt/system/calc/server.js"
pidfile="/run/${RC_SVCNAME}.pid"
name="dyno"
description="dyno server"

depend() { need net; }
